# picoclaw-launcher.ini
[program:picoclaw-launcher]
user=root
environment=GIN_MODE="release",LD_LIBRARY_PATH="",HOME="/root",PATH="/sbin:/bin:/usr/local/bin:/root/.local/bin:/usr/local/picoclaw",PICO_VER_URL="https://lib10.cn/download/picoclaw-ver.json"
directory=/usr/local/picoclaw
command=./picoclaw-launcher -https=:18800 -http=:18801 -discovery=:18801
autostart=true
autorestart=true
startsecs=3
startretries=3
autorestart=true
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/vpot/var/log/picoclaw-launcher.logfile.stdout
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/vpot/var/log/picoclaw-launcher.logfile.stderr
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
